# deploy/guidebot-voice.service
[Unit]
Description=Guidebot Qwen realtime voice
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Guidebot
EnvironmentFile=/etc/guidebot/voice.env
Environment=PYTHONUTF8=1
Environment=LANG=C.UTF-8
ExecStart=/home/pi/Guidebot/.venv/bin/guidebot voice-qwen --voice Tina
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
